# sources/systemd/scx.service.in
[Unit]
Description=sched-ext BPF CPU Scheduler (@SCHEDULER@)
After=sysinit.target local-fs.target
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/@SCHEDULER@
TimeoutStopSec=10

[Install]
WantedBy=sysinit.target
